# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
childlogdir=/var/log/supervisor

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:isolate-panel]
command=/bin/sh -c "PORT=8080 CONFIG_PATH=/app/configs/config.yaml DATABASE_PATH=/app/data/isolate-panel.db LOG_LEVEL=info APP_ENV=production exec /usr/local/bin/isolate-panel"
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/isolate-panel/panel.err.log
stdout_logfile=/var/log/isolate-panel/panel.out.log
user=root
priority=10

[program:xray]
command=/app/data/cores/xray/xray -config /app/data/cores/xray/config.json
directory=/app/data/cores/xray
autostart=false
autorestart=true
stderr_logfile=/var/log/isolate-panel/xray.err.log
stdout_logfile=/var/log/isolate-panel/xray.out.log
user=root
priority=20
startsecs=5
startretries=3
stopwaitsecs=10

[program:mihomo]
command=/app/data/cores/mihomo/mihomo -d /app/data/cores/mihomo -f /app/data/cores/mihomo/config.yaml
directory=/app/data/cores/mihomo
autostart=false
autorestart=true
stderr_logfile=/var/log/isolate-panel/mihomo.err.log
stdout_logfile=/var/log/isolate-panel/mihomo.out.log
user=root
priority=20
startsecs=5
startretries=3
stopwaitsecs=10

[program:singbox]
command=/app/data/cores/singbox/sing-box run -c /app/data/cores/singbox/config.json
directory=/app/data/cores/singbox
autostart=false
autorestart=true
stderr_logfile=/var/log/isolate-panel/singbox.err.log
stdout_logfile=/var/log/isolate-panel/singbox.out.log
user=root
priority=20
startsecs=5
startretries=3
stopwaitsecs=10
