# docker/supervisord.fullstack.conf
; NOTE: Credentials are read from OS environment variables. Set JWT_SECRET and ADMIN_PASSWORD before starting.
;
; Supervisord configuration for full-stack integration tests
; This configuration AUTOSTARTS all cores for full-stack testing
;

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=root:root

[inet_http_server]
port=127.0.0.1:9001

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info
childlogdir=/var/log/supervisor
user=root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================================================
; Main Application
; ============================================================================

[program:isolate-panel]
command=/usr/local/bin/isolate-panel
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
exitcodes=0,2
priority=10
user=isolate
environment=PORT="8080",DATABASE_PATH="/app/data/isolate-panel.db",LOG_LEVEL="debug",APP_ENV="development",JWT_SECRET="%(ENV_JWT_SECRET)s",ISOLATE_JWT_SECRET="%(ENV_JWT_SECRET)s",ADMIN_USERNAME="admin",ADMIN_PASSWORD="%(ENV_ADMIN_PASSWORD)s",MONITORING_MODE="full",TZ="UTC"
stdout_logfile=/var/log/supervisor/isolate-panel.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/isolate-panel.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3

; ============================================================================
; Xray Core - AUTOSTART enabled for full-stack tests
; ============================================================================

[program:xray]
command=/usr/local/bin/xray run -config /shared-configs/xray.json
directory=/app/data/cores/xray
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
exitcodes=0,2
priority=20
user=isolate
stdout_logfile=/var/log/supervisor/xray.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/xray.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3

; ============================================================================
; Sing-box Core - AUTOSTART enabled for full-stack tests
; ============================================================================

[program:singbox]
command=/usr/local/bin/sing-box run -c /shared-configs/singbox.json --disable-color
directory=/app/data/cores/singbox
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
exitcodes=0,2
priority=20
user=isolate
stdout_logfile=/var/log/supervisor/singbox.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/singbox.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3

; ============================================================================
; Mihomo Core - AUTOSTART enabled for full-stack tests
; ============================================================================

[program:mihomo]
command=/usr/local/bin/mihomo -d /shared-configs -f /shared-configs/mihomo.yaml
directory=/app/data/cores/mihomo
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
exitcodes=0,2
priority=20
user=isolate
stdout_logfile=/var/log/supervisor/mihomo.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/mihomo.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3

; ============================================================================
; Process Groups — Core Isolation
; ============================================================================

[group:xray]
programs=xray
priority=100

[group:mihomo]
programs=mihomo
priority=110

[group:singbox]
programs=singbox
priority=120
